# systemd/earnings-gap-trader-worker.service
# Systemd service file for Earnings Gap Trading System - Celery Worker
# Location: /etc/systemd/system/earnings-gap-trader-worker.service

[Unit]
Description=Earnings Gap Trading System - Background Worker
Documentation=https://github.com/your-org/earnings-gap-trader
After=network.target postgresql.service redis.service rabbitmq-server.service
Wants=network.target
Requires=postgresql.service redis.service rabbitmq-server.service

[Service]
Type=exec
User=trader
Group=trader
WorkingDirectory=/opt/earnings_gap_trader/app
Environment=PATH=/opt/earnings_gap_trader/venv/bin
EnvironmentFile=/opt/earnings_gap_trader/config/.env

# Celery worker process
ExecStart=/opt/earnings_gap_trader/venv/bin/celery \
    --app=core.celery_app \
    worker \
    --loglevel=info \
    --concurrency=2 \
    --hostname=worker@%%h \
    --queues=trading,signals,notifications,maintenance

# Graceful shutdown
ExecStop=/bin/kill -TERM $MAINPID

# Restart configuration
Restart=always
RestartSec=10
StartLimitBurst=3
StartLimitInterval=300

# Resource limits
LimitNOFILE=65536
LimitNPROC=2048
MemoryLimit=1G
CPUQuota=50%

# Security measures
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/earnings_gap_trader /var/log/earnings_gap_trader /opt/backups/earnings_gap_trader
PrivateDevices=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictRealtime=true
RestrictNamespaces=true

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=60

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=earnings-gap-trader-worker

# Working directory permissions
UMask=0027

[Install]
WantedBy=multi-user.target